# com.davidjesse.running-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.davidjesse.running-dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/davidjesse/running-dashboard/run_sync.sh</string>
    </array>

    <!-- Run at 6:00 AM local time.
         launchd will fire once on the next wake if the Mac was asleep at
         the scheduled time — no extra configuration needed. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Provide the minimal environment launchd strips out. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/davidjesse</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Don't run on install, only on schedule. -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
