# com.davidjesse.sync-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.davidjesse.sync-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/davidjesse/running-dashboard/sync_server.py</string>
    </array>

    <!-- Start on login and keep it running permanently in the background. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Provide the minimal environment launchd strips out. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/davidjesse</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Process stdout/stderr (server start/stop, tracebacks) — separate from
         the sync-trigger log the server writes itself. -->
    <key>StandardOutPath</key>
    <string>/Users/davidjesse/running-dashboard-sync-server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/davidjesse/running-dashboard-sync-server.log</string>
</dict>
</plist>
